# system/rootfs/var/lib/dpkg/info/mosquitto.list
/.
/etc
/etc/init.d
/etc/init.d/mosquitto
/etc/logrotate.d
/etc/logrotate.d/mosquitto
/etc/mosquitto
/etc/mosquitto/ca_certificates
/etc/mosquitto/ca_certificates/README
/etc/mosquitto/certs
/etc/mosquitto/certs/README
/etc/mosquitto/conf.d
/etc/mosquitto/conf.d/README
/etc/mosquitto/mosquitto.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/mosquitto.service
/usr
/usr/bin
/usr/bin/mosquitto_passwd
/usr/sbin
/usr/sbin/mosquitto
/usr/share
/usr/share/doc
/usr/share/doc/mosquitto
/usr/share/doc/mosquitto/README.Debian
/usr/share/doc/mosquitto/changelog.Debian.gz
/usr/share/doc/mosquitto/changelog.gz
/usr/share/doc/mosquitto/copyright
/usr/share/doc/mosquitto/examples
/usr/share/doc/mosquitto/examples/aclfile.example
/usr/share/doc/mosquitto/examples/mosquitto.conf.gz
/usr/share/doc/mosquitto/examples/pskfile.example
/usr/share/doc/mosquitto/examples/pwfile.example
/usr/share/doc/mosquitto/readme.md
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mosquitto_passwd.1.gz
/usr/share/man/man5
/usr/share/man/man5/mosquitto.conf.5.gz
/usr/share/man/man7
/usr/share/man/man7/mosquitto-tls.7.gz
/usr/share/man/man7/mqtt.7.gz
/usr/share/man/man8
/usr/share/man/man8/mosquitto.8.gz
/var
/var/lib
/var/lib/mosquitto
/var/log
/var/log/mosquitto
